# hhh.conf
[program:hhh]
command=python3 /home/curator2/hhh/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/hhh/test.err.log
stdout_logfile=/home/curator2/hhh/test.out.log